# build/usb/includes.chroot/etc/systemd/system/ssh.service.d/rigos-observe.conf
[Unit]
After=rigos-recovery-access.service
Wants=rigos-remote-access-observe.service
