# build/usb/includes.chroot/etc/systemd/system/rigos-firstboot.service.d/2009-console-theme.conf
[Service]
Environment=RIGOS_WHIPTAIL=/usr/lib/rigos/rigos-firstboot-whiptail
Environment="RIGOS_FIRSTBOOT_BACKTITLE=RIGOS SYSTEM CONFIGURATION // LOCAL NODE SETUP // OFFLINE AUTHORITY"
