# build/usb/includes.chroot/etc/systemd/system/rigos-state.service
[Unit]
Description=Mount RIGOS persistent state
DefaultDependencies=no
After=systemd-udevd.service local-fs-pre.target
Before=local-fs.target rigos-state-ready.service

[Service]
Type=oneshot
ExecStartPre=/usr/bin/systemd-tmpfiles --create /usr/lib/tmpfiles.d/rigos.conf
ExecStart=/usr/lib/rigos/rigos-state-init
RemainAfterExit=yes

[Install]
WantedBy=local-fs.target
