# build/usb/includes.chroot/etc/systemd/system/rigos-state-ready.service
[Unit]
Description=Verify RIGOS persistent state readiness
DefaultDependencies=no
After=rigos-state.service
Requires=rigos-state.service
Before=local-fs.target rigos-profile-apply.service rigos-firstboot.service rigos-hugepages.service rigos-miner.service

[Service]
Type=oneshot
ExecStart=/usr/lib/rigos/rigos-state-ready
RemainAfterExit=yes

[Install]
WantedBy=local-fs.target
