# build/usb/includes.chroot/etc/systemd/system/rigos-state.service
[Unit]
Description=Mount RIGOS persistent state
DefaultDependencies=no
After=systemd-udev-settle.service local-fs-pre.target
Before=local-fs.target rigos-firstboot.service rigos-miner.service

[Service]
Type=oneshot
ExecStart=/usr/lib/rigos/rigos-state-init
RemainAfterExit=yes

[Install]
WantedBy=local-fs.target
