# build/usb/includes.chroot/etc/systemd/system/rigos-state-ready.service
[Unit]
Description=Verify RIGOS persistent state readiness
After=rigos-state.service rigos-recovery-access.service
Requires=rigos-state.service
Before=rigos-ssh-hostkeys.service rigos-profile-apply.service rigos-firstboot.service rigos-hugepages.service rigos-miner.service

[Service]
Type=oneshot
ExecStart=/usr/lib/rigos/rigos-state-ready
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
