# build/usb/includes.chroot/etc/systemd/system/rigos-firstboot.service
[Unit]
Description=RIGOS first-boot configuration
After=network-online.target rigos-state-ready.service rigos-profile-apply.service
Wants=network-online.target
Requires=rigos-state-ready.service
Before=getty@tty1.service

[Service]
Type=oneshot
ExecCondition=/usr/lib/rigos/rigos-config needs-activation
ExecStart=/usr/local/sbin/rigos-firstboot-seeded
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVTDisallocate=yes

[Install]
WantedBy=multi-user.target
