# build/usb/includes.chroot/etc/systemd/system/rigos-miner.service.d/randomx-msr.conf
[Unit]
Wants=rigos-randomx-msr.service
After=rigos-randomx-msr.service
