# build/usb/includes.chroot/etc/systemd/system/rigos-miner.service.d/stability.conf
[Unit]
StartLimitIntervalSec=10min
StartLimitBurst=5

[Service]
Restart=on-failure
RestartSec=15s
TimeoutStartSec=2min
TimeoutStopSec=30s
KillSignal=SIGINT
FinalKillSignal=SIGKILL
SendSIGKILL=yes
SuccessExitStatus=SIGINT SIGTERM
OOMPolicy=stop
OOMScoreAdjust=-500
TasksMax=64
LimitNOFILE=4096
