# build/usb/includes.chroot/etc/systemd/system/getty@tty1.service.d/rigos-firstboot.conf
[Unit]
Wants=rigos-firstboot.service
After=rigos-firstboot.service
